# rpi/services/mia-gpio-worker.service
[Unit]
Description=MIA GPIO Worker
After=network.target zmq-broker.service
Wants=zmq-broker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-servis/rpi
ExecStart=/usr/bin/python3 /opt/ai-servis/rpi/hardware/gpio_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
