# rpi/services/mia-api.service
[Unit]
Description=MIA FastAPI Server for Raspberry Pi
After=network.target zmq-broker.service
Wants=zmq-broker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mia/rpi
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
