# rpi/services/mia-obd-worker.service
[Unit]
Description=MIA OBD Worker - ELM327 Simulator Service
After=network.target zmq-broker.service mia-serial-bridge.service
Wants=zmq-broker.service mia-serial-bridge.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-servis/rpi
ExecStart=/usr/bin/python3 /opt/ai-servis/rpi/services/obd_worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
